# systemd/aws-key-rotate.timer
[Unit]
Description=Rotate AWS Keys

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
