# systemd/aws-key-rotate.service
[Unit]
Description=Rotate AWS Keys

[Service]
Type=oneshot
ExecStart=/home/jkomoroski/.local/bin/aws-key-rotate
