# systemd/aws-key-rotate.service
[Unit]
Description=Rotate AWS Keys

[Service]
Type=oneshot
ExecStart=%h/.local/bin/aws-rotate-keys
